# keys.timer
[Unit]
Description=Update SSH keys once a day

[Timer]
OnCalendar=daily
Unit=keys.service

[Install]
WantedBy=timers.target
